# crontab
*/5 * * * * php /var/www/shop.mambasoft.it/cron.php
